# services/maraschino/run
#!/bin/bash
/sbin/setuser abc python /app/maraschino/Maraschino.py --datadir /config

