# root/etc/services.d/maraschino/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc python /app/maraschino/Maraschino.py \
	--datadir /config
